# systemd/hf-tec@.service
[Unit]
Description=HF PRN beacon recorder for ka9q-radio (%i)
Documentation=https://github.com/HamSCI/hf-tec
After=network-online.target radiod@%i.service
Wants=network-online.target

[Service]
Type=notify
User=hftec
Group=hftec
ExecStart=/usr/local/bin/hf-tec daemon --instance %i
Restart=on-failure
RestartSec=15
TimeoutStopSec=30

# Resilience backstops (parity with codar-sounder / hfdl-recorder): arm the
# watchdog the daemon already feeds (WATCHDOG=1 every WatchdogSec/2), and cap
# memory so a runaway buffer can't OOM the host — it gets restarted instead.
WatchdogSec=180
MemoryMax=512M
MemorySwapMax=0

# Logging — daemon process log goes to systemd journal; per-frame data
# goes to /var/lib/hf-tec/<radiod_id>/.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hf-tec@%i

# Environment — sigmond's coordination.env supplies per-radiod chain delays
# and runtime log-level overrides.
EnvironmentFile=-/etc/sigmond/coordination.env
Environment=PYTHONUNBUFFERED=1

# Sandboxing — standard recorder profile.
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/hf-tec /var/log/hf-tec /var/lib/sigmond
RuntimeDirectory=hf-tec
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
